# cron/crontab
* * * * * php /var/www/html/crear.php
# An empty line
